# profile.d/outlook.sh
function owa(){
nohup chromium --app=https://outlook.autoscout24.com/owa/ >& /dev/null &
}
